# deploy/systemd/datalevin-docs.service
[Unit]
Description=Datalevin documentation site
After=network.target

[Service]
Type=simple
User=datalevin
Group=datalevin
WorkingDirectory=/opt/datalevin-docs
EnvironmentFile=-/etc/datalevin-docs/datalevin-docs.env
ExecStart=/usr/bin/java -Xms256m -Xmx384m -XX:+UseG1GC -XX:MaxGCPauseMillis=200 -XX:+UseStringDeduplication -XX:MaxMetaspaceSize=128m -jar /opt/datalevin-docs/datalevin-docs-standalone.jar
MemoryHigh=800M
MemoryMax=900M
OOMPolicy=stop
Restart=on-failure
RestartSec=5
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
